# scripts/cron/backup.cron
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
0 0 * * *       run-parts --report --arg=daily --arg=8 /home/vcm/cron/jobs >> /home/vcm/cron/logs/backup.log
0 0 * * 0	run-parts --report --arg=weekly --arg=5 /home/vcm/cron/jobs >> /home/vcm/cron/logs/backup.log
0 0 1 * *	run-parts --report --arg=monthly --arg=13 /home/vcm/cron/jobs >> /home/vcm/cron/logs/backup.log

